# smx-dashboard.service
[Unit]
Description=Gunicorn instance to serve StepManiaX B2B Dashboard
After=network.target

[Service]
User=jules
Group=jules
WorkingDirectory=/app
Environment="PATH=/app/venv/bin"
Environment="SECRET_KEY=replace-with-secure-key"
Environment="SMX_API_KEY=replace-with-secure-key"
ExecStart=/app/venv/bin/gunicorn --config gunicorn_config.py app:app

[Install]
WantedBy=multi-user.target
